# puppet/modules/etcd/templates/etcd.service.erb
[Unit]
Description=etcd
Documentation=https://github.com/coreos

[Service]
Type=notify
ExecStart=/usr/bin/etcd \
  --name <%= @hostname %> \
  --cert-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes.pem \
  --key-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes-key.pem \
  --peer-cert-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes.pem \
  --peer-key-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes-key.pem \
  --trusted-ca-file=<%= scope.lookupvar('kube_common::conf_dir') %>/ca.pem \
  --peer-trusted-ca-file=<%= scope.lookupvar('kube_common::conf_dir') %>/ca.pem \
  --peer-client-cert-auth \
  --client-cert-auth \
  --initial-advertise-peer-urls https://<%= @ipaddress %>:2380 \
  --listen-peer-urls https://<%= @ipaddress %>:2380 \
  --listen-client-urls https://<%= @ipaddress %>:2379,https://127.0.0.1:2379 \
  --advertise-client-urls https://<%= @ipaddress %>:2379 \
  --initial-cluster-token etcd-cluster-0 \
  --initial-cluster <%= @network_info['controllers'].map { |name, info| "#{name.split('.').first}=https://#{info['private_ip']}:2380" }.join(',')  %> \
  --initial-cluster-state new \
  --data-dir=/var/lib/etcd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target