# puppet/modules/kube_scheduler/templates/kube-scheduler.service.erb
Description=Kubernetes Scheduler
Documentation=https://github.com/kubernetes/kubernetes

[Service]
ExecStart=/usr/bin/kube-scheduler \
  --config=<%= scope.lookupvar('kube_common::conf_dir') %>/kube-scheduler.yaml \
  --v=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
