# puppet/modules/kube_apiserver/templates/kube-apiserver.service.erb
[Unit]
Description=Kubernetes API Server
Documentation=https://github.com/kubernetes/kubernetes

[Service]
ExecStart=/usr/bin/kube-apiserver \
  --advertise-address=<%= @ipaddress %> \
  --allow-privileged=true \
  --apiserver-count=3 \
  --audit-log-maxage=30 \
  --audit-log-maxbackup=3 \
  --audit-log-maxsize=100 \
  --audit-log-path=/var/log/audit.log \
  --authorization-mode=Node,RBAC \
  --bind-address=0.0.0.0 \
  --client-ca-file=<%= scope.lookupvar('kube_common::conf_dir') %>/ca.pem \
  --enable-admission-plugins=NamespaceLifecycle,NodeRestriction,LimitRanger,ServiceAccount,DefaultStorageClass,ResourceQuota \
  --etcd-cafile=<%= scope.lookupvar('kube_common::conf_dir') %>/ca.pem \
  --etcd-certfile=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes.pem \
  --etcd-keyfile=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes-key.pem \
  --etcd-servers=<%= @network_info['controllers'].map { |name, info| "https://#{info['private_ip']}:2379" }.join(',')  %> \
  --event-ttl=1h \
  --encryption-provider-config=<%= scope.lookupvar('kube_common::conf_dir') %>/encryption-config.yaml \
  --kubelet-certificate-authority=<%= scope.lookupvar('kube_common::conf_dir') %>/ca.pem \
  --kubelet-client-certificate=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes.pem \
  --kubelet-client-key=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes-key.pem \
  --kubelet-https=true \
  --runtime-config='api/all=true' \
  --service-account-key-file=<%= scope.lookupvar('kube_common::conf_dir') %>/service-account.pem \
  --service-cluster-ip-range=<%= @network_info['service_cluster_ip_range'] %> \
  --service-node-port-range=30000-32767 \
  --tls-cert-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes.pem \
  --tls-private-key-file=<%= scope.lookupvar('kube_common::conf_dir') %>/kubernetes-key.pem \
  --v=2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
